# ockam-server/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/ockam/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ockam-node]
; The Ockam transport listener that the outside world reaches on :14000.
; --foreground keeps the process alive so supervisord can monitor it.
command=/usr/local/bin/ockam node create %(ENV_OCKAM_NODE_NAME)s --tcp-listener-address %(ENV_OCKAM_NODE_TRANSPORT)s --foreground
user=ockam
environment=OCKAM_HOME="/var/lib/ockam-server",HOME="/var/lib/ockam-server"
autostart=true
autorestart=true
startretries=10
startsecs=2
priority=10
stdout_logfile=/var/log/ockam/node.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/ockam/node.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:ockam-controller]
; The control-plane FastAPI app, bound to 127.0.0.1:8080 inside the
; container (NOT exposed to the outside).
command=/opt/venv/bin/python -m ockam_controller --bind 127.0.0.1:8080 --state %(ENV_OCKAM_CONTROLLER_STATE)s --log-level info
user=ockam
environment=OCKAM_HOME="/var/lib/ockam-server",HOME="/var/lib/ockam-server",OCKAM_NODE_NAME="%(ENV_OCKAM_NODE_NAME)s",OCKAM_NODE_TRANSPORT="%(ENV_OCKAM_NODE_TRANSPORT)s",OCKAM_BINARY="%(ENV_OCKAM_BINARY)s",OCKAM_CONTROLLER_ADMIN_IDENTIFIERS="%(ENV_OCKAM_CONTROLLER_ADMIN_IDENTIFIERS)s",OCKAM_CONTROLLER_TRUST_ALL="%(ENV_OCKAM_CONTROLLER_TRUST_ALL)s"
autostart=true
autorestart=true
startretries=10
startsecs=3
; Wait for ockam-node to be up before starting the controller — the
; controller's reconcile() at startup will try to talk to the node.
priority=20
stdout_logfile=/var/log/ockam/controller.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/ockam/controller.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
